# scripts/helifx.service
[Unit]
Description=Helicopter FX System
After=network.target sound.target alsa-restore.service multi-user.target
DefaultDependencies=no

[Service]
Type=simple
User={{USER}}
Group={{USER}}
SupplementaryGroups=audio gpio
WorkingDirectory={{INSTALL_DIR}}
ExecStart={{INSTALL_DIR}}/helifx {{CONFIG_PATH}}
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=false

# Environment
Environment="HOME={{USER_HOME}}"
Environment="PULSE_RUNTIME_PATH=/run/user/{{USER_UID}}/pulse"
Environment="XDG_RUNTIME_DIR=/run/user/{{USER_UID}}"

[Install]
WantedBy=multi-user.target
