# scripts/helifx-audio.service
[Unit]
Description=HeliFX Audio Level Setup
Documentation=https://github.com/yourusername/helifx
DefaultDependencies=no
After=multi-user.target sound.target graphical.target systemd-user-sessions.service
Before=helifx.service
Wants=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/helifx-audio-setup --verbose
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
