# scripts/helifx.service
[Unit]
Description=Helicopter FX System
After=network.target sound.target alsa-restore.service multi-user.target
DefaultDependencies=no

[Service]
Type=simple
User=helisfx
Group=helisfx
SupplementaryGroups=audio gpio
WorkingDirectory=/home/pi/helifx
ExecStart=/home/pi/helifx/helifx /home/pi/helifx/config.yaml
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=false

# Environment
Environment="HOME=/home/helisfx"
Environment="PULSE_RUNTIME_PATH=/run/user/1000/pulse"
Environment="XDG_RUNTIME_DIR=/run/user/1000"

[Install]
WantedBy=multi-user.target
